# supervisor.d/feed-dashboard.consume.conf
[supervisord]
logfile=/project/var/log/supervisord.log
loglevel=error

[program:messenger-consume-feed-dashboard]
command=php /project/bin/console messenger:consume async_dashboard_feeding --time-limit=3600
numprocs=3
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d
